# docker/config/supervisor/cups.conf
[program:cups]
command=/usr/sbin/cupsd -f
